# files/profile.d/expiry.sh.new
#!/bin/bash
echo -en "\033[1;31m"
/usr/bin/expiry -c
echo -en "\033[0m"
