# newconfs/cron.d/rkhunter.new
30 * * * * /bin/test -x /usr/bin/rkhunter && /usr/bin/rkhunter --cronjob --rwo
