# files/udisks2.service
[Service]
MemoryDenyWriteExecute=true
ProtectClock=true
ProtectHostname=true
LockPersonality=true
SystemCallArchitectures=native
